# deploy/fundtrader-backfill.timer


[Unit]
Description=Fund detail snapshots backfill (every 6 hours)

[Timer]
OnCalendar=*-*-* 02,08,14,20:00:00
Persistent=true

[Install]
WantedBy=timers.target