# deploy/fundtrader-webhook.service
[Unit]
Description=FundTrader Gitee Webhook Listener
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/fundtrader
Environment=WEBHOOK_SECRET=fundtrader_webhook_2024
Environment=WEBHOOK_PORT=9100
Environment=GIT_REMOTE=gitee
Environment=BRANCH=master
ExecStart=/usr/bin/python3 /opt/fundtrader/scripts/webhook-listener.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
